# pictureframe.service
[Unit]
Description=pictureframe
After=multi-user.target

[Service]
ExecStartPre=/bin/sleep 30
TimeoutStopSec=2
Type=simple
Restart=always
ExecStart=/opt/pictureframe.py
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
